# meta-phire/recipes-core/set-hostname/files/set-hostname.service
[Unit]
Description=Set hostname from file
Before=systemd-networkd.service
After=data.mount

[Service]
Type=oneshot
ExecStart=/usr/bin/set-hostname

[Install]
WantedBy=multi-user.target
